# kiro-wecom-bridge.service
[Unit]
Description=Kiro WeCom Bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/mnt/i/workspace/kiro-wecom-bridge
ExecStart=/mnt/i/workspace/kiro-wecom-bridge/.venv/bin/python3 main.py
Restart=always
RestartSec=5
Environment=PATH=/root/.local/bin:/root/.cargo/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
